# etc/supervisor.d/plackup.ini
[program:plackup]
command=plackup --env %(ENV_PLACK_ENV)s --access-log /dev/null --port 9000 /app.psgi
user=otrs
autostart=%(ENV_START_PLACKUP)s
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0